# pkg/ultimate-gaming-server/usr/lib/systemd/user/sunshine-watchdog.service
[Unit]
Description=Sunshine idle connection watchdog
# Démarre après Sunshine, s'arrête automatiquement si Sunshine s'arrête
BindsTo=sunshine.service
After=sunshine.service

[Service]
Type=simple
ExecStart=/usr/bin/sunshine-watchdog
Restart=no
