# systemd/sunshine.service
[Unit]
Description=Sunshine Game Stream Server
After=graphical-session.target xdg-desktop-portal-hyprland.service xorg-stream.service
Requires=xorg-stream.service
Conflicts=sunshine.socket
Wants=sunshine-watchdog.service

[Service]
# DISPLAY est lu depuis /tmp/xorg-stream-info (écrit par start-xorg-stream)
Environment=XAUTHORITY=/tmp/xauth-stream
# Supprime WAYLAND_DISPLAY pour que Steam/Lutris/jeux utilisent X11 (DISPLAY=:2)
# et soient capturés par NvFBC. Sans ça ils s'ouvrent sur Hyprland et NvFBC voit noir.
UnsetEnvironment=WAYLAND_DISPLAY
ExecStartPre=-/usr/bin/pkill -x sunshine
ExecStartPre=/bin/sh -c '\
    for i in $(seq 1 20); do \
        DISP=$(grep "^display=" /tmp/xorg-stream-info 2>/dev/null | cut -d= -f2); \
        [ -n "$DISP" ] && systemctl --user set-environment DISPLAY="$DISP" && exit 0; \
        sleep 0.5; \
    done; exit 1'
ExecStart=/usr/bin/start-sunshine

ExecStopPost=systemctl --user --no-block stop sunshine-watchdog.service
ExecStopPost=systemctl --user --no-block start sunshine.socket

Restart=no

[Install]
WantedBy=default.target
