# systemd/xorg-stream-system.service
[Unit]
Description=X11 headless stream server (NvFBC)
# Démarre après logind (nécessaire pour la gestion des VT) mais AVANT le display manager
# (SDDM/Hyprland). Xorg démarre ainsi sur un VT libre sans aucun compositor actif.
After=systemd-logind.service
Before=display-manager.service

[Service]
Type=simple
User=acidtracks
Environment=XAUTHORITY=/tmp/xauth-stream
ExecStart=/usr/bin/start-xorg-stream
ExecStartPost=/bin/sh -c '\
    for i in $(seq 1 30); do \
        DISP=$(grep "^display=" /tmp/xorg-stream-info 2>/dev/null | cut -d= -f2); \
        [ -n "$DISP" ] && XAUTHORITY=/tmp/xauth-stream xauth add "$DISP" . $(mcookie) && exit 0; \
        sleep 0.5; \
    done; exit 1'
ExecStopPost=/bin/sh -c 'rm -f /tmp/xauth-stream /tmp/xorg-stream-info'
Restart=no

[Install]
WantedBy=graphical.target
